# poolapi.service
[Unit]
Description=TeamB_Pool_Project_API
After=network.target

[Service]
WorkingDirectory=/home/terraform
ExecStart=/usr/bin/node /terraform/api/TeamB_Pool_Project_API/index.js
Restart=on-failure

[Install]
WantedBy=multi-user.target