# poolapi.service
[Unit]
Description=TeamB_Pool_Project_API
After=network.target

[Service]
WorkingDirectory=/home/api
ExecStart=/usr/bin/node /home/api/TeamB_Pool_Project_API/index.js
Restart=on-failure

[Install]
WantedBy=multi-user.target